# packages/portalocker/extra/dcos-portalocker.windows.service
[Unit]
Description=Properly install the Python portalocker dependent library pywin32.
ConditionPathExists=!c:\opt\mesosphere\etc\pywin32.initialized
[Service]
Type=oneshot
StandardOutput=journal+console
StandardError=journal+console
User=dcos_portalocker
EnvironmentFile=c:\opt\mesosphere\environment.export.ps1
ExecStart=pip install --upgrade --force-reinstall pywin32
ExecStartPost=-cmd.exe /c mkdir  c:\opt\mesosphere\etc\pywin32.initialized
